# roles/install_klipper/files/klipper-pru.path
#Systemd path file for klipper-pru
[Unit]
Description=klipper pru PATH
PartOf=klipper_pru.service

[Path]
PathExists=/sys/class/remoteproc/remoteproc2/state
PathExists=/sys/class/remoteproc/remoteproc1/state
