# roles/install_klipper/files/remoteproc1.path
[Unit]
Description=Remoteproc1 availability monitoring

[Path]
PathExists=/dev/remoteproc/
Unit=klipper_pru.service

[Install]
WantedBy=multi-user.target
